# pods.service
[Unit]
Description=Pods MCP Server
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/pods
ExecStart=PYTHON_PLACEHOLDER server.py --http
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
